# launchd/com.user.codex-token-dashboard-recalc.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.user.codex-token-dashboard-recalc</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/ABSOLUTE/PATH/TO/codex-token-usage-dashboard/scripts/run_local.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/ABSOLUTE/PATH/TO/codex-token-usage-dashboard</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/com.user.codex-token-dashboard-recalc.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/com.user.codex-token-dashboard-recalc.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>CODEX_USAGE_SERVER_HOST</key>
      <string>127.0.0.1</string>
      <key>CODEX_USAGE_SERVER_PORT</key>
      <string>8765</string>
      <key>CODEX_USAGE_SESSIONS_ROOT</key>
      <string>/ABSOLUTE/PATH/TO/.codex/sessions</string>
      <key>CODEX_USAGE_CLAUDE_PROJECTS_ROOT</key>
      <string>/ABSOLUTE/PATH/TO/.claude/projects</string>
    </dict>
  </dict>
</plist>
